# net_config.cron
* * * * * root /usr/scripts/net-config.sh
